# examples/launchd/com.example.open-wispr-nightly-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.example.open-wispr-nightly-review</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/ABSOLUTE/PATH/TO/scripts/openwispr_nightly_review.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/REPO</string>
    <key>StandardOutPath</key>
    <string>/tmp/open-wispr-nightly-review.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/open-wispr-nightly-review.log</string>
  </dict>
</plist>

<!-- #huanyuan -->
